# rootfs/etc/profile.d/02-nvm.sh
#!/usr/bin/env bash

. "$NVM_DIR/nvm.sh"

nvm use ${DEFAULT_NODE_VERSION:---lts}
